# axiscope.service
[Unit]
Description=AxisScope - Tool Alignment Interface for Klipper
After=network.target moonraker.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=$USER
WorkingDirectory=$HOME/axiscope
ExecStart=$HOME/axiscope/axiscope-env/bin/python3 app.py
Environment="PATH=$HOME/axiscope/axiscope-env/bin"
Environment="PYTHONUNBUFFERED=1"
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
